# soc-hub.service
[Unit]
Description=CD-Startpage — Service Dashboard
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=$HOME/soc-hub
ExecStart=/usr/bin/python3 $HOME/soc-hub/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
